# rootfs/etc/s6-overlay/s6-rc.d/pulseaudio/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start PulseAudio service
# ==============================================================================

# Set env
export PULSE_STATE_PATH="/data/states"
export LD_PRELOAD="/usr/local/lib/libjemalloc.so.2"

exec s6-notifyoncheck -d -s 300 -w 300 -n 0 \
     pulseaudio --system --disallow-exit --exit-idle-time=-1 \
                --disable-shm "${pulse_args[@]}"